# launchd/Library/LaunchAgents/local.safari.bookmarks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <!--
        This launch agent imports new bookmarks into Safari when the bookmarks backup is updated.
        -->
        <key>Label</key>
        <string>local.safari.bookmarks</string>
        <!-- TODO: set WorkingDirectory via `defaults` instead of hardcoding -->
        <key>WorkingDirectory</key>
        <string>/Users/tmz</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/bin/python3</string>
            <string>.dotfiles/scripts/import_safari_bookmarks.py</string>
            <string>.backup/bookmarks.html</string>
        </array>
        <key>WatchPaths</key>
        <array>
            <string>.backup/bookmarks.html</string>
        </array>
    </dict>
</plist>
