# launchd/Library/LaunchAgents/local.safari.bookmarks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>local.safari.bookmarks</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/bin/python3</string>
            <string>/Users/tmz/.dotfiles/scripts/import_safari_bookmarks.py</string>
            <string>/Users/tmz/backup/bookmarks.html</string>
        </array>
        <key>WatchPaths</key>
        <array>
            <string>/Users/tmz/backup/bookmarks.html</string>
        </array>
    </dict>
</plist>
